# flask_app/supervisord.conf
[supervisord]
nodaemon=true

[program:nginx]
command=/usr/sbin/nginx -g "daemon off;"
autostart=true
autorestart=true
priority=1

[program:gunicorn]
command=/usr/local/bin/gunicorn app:app --bind 0.0.0.0:8080
autostart=true
autorestart=true
priority=2